# cron/crontab
*/5 * * * * cd /checklab/src && /usr/local/bin/python3 check_email_box.py >> /var/log/cron.log 2>&1
